# templates/systemd/atop-rotate.service.j2
[Unit]
Description=Restart atop daemon to rotate logs

[Service]
Type=oneshot
ExecStart=/bin/systemctl try-restart atop.service
